# supervisord.fluxbox.conf
[program:fluxbox]
command = /usr/local/bin/run-as-x11-user
	/usr/bin/fluxbox
	-display
	:"%(ENV_NUMBER_DISPLAY)s"
startsecs = 10
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
